# kbd-backlight-autosave.service
[Unit]
Description=Periodically save keyboard backlight state

[Service]
Type=simple
ExecStart=/usr/local/bin/led save
